# systemd/reef-n1.service
[Unit]
Description=Reef Citadel — N1 Darwin (pure fitness, no gauntlet)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/reef_citadel
Environment=PYTHONPATH=%h/reef_citadel
Environment=REEF_COLONY=n1_darwin
Environment=REEF_DATA_ROOT=%h/apex_data
Environment=REEF_STATE_ROOT=%h/reef_citadel/.state
ExecStart=/usr/bin/python3 -m colonies.n1_darwin.run
Restart=always
RestartSec=30
Nice=10
CPUQuota=40%
MemoryMax=2G
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
